# scripts/service/dbus-daemon/run
#!/bin/sh

set -eu

rm -rf /run/dbus/
mkdir -p /run/dbus/

exec 2>&1
exec dbus-daemon --system --nofork --nopidfile
